# crontab/root
# min hour day month weekday command
*/2 * * * * cd /app && uv run main.py > /dev/stdout