# crontab/root
# min hour day month weekday command
0 * * * * cd /app && uv run main.py > /dev/stdout